# shard58/harbor-shard58.service
[Unit]
Description=Harbor Shard 58 - bako-biib P2P Chat
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/mdupont/introspector/shard58/harbor
Environment="PATH=/nix/var/nix/profiles/default/bin:/usr/bin:/bin"
ExecStart=/usr/bin/bash -c 'cd /home/mdupont/introspector/shard58/harbor && /nix/var/nix/profiles/default/bin/nix develop -c npm run tauri dev'

# Filesystem protection
ProtectSystem=strict
ReadWritePaths=/tmp /home/mdupont/introspector/shard58/harbor

# Security
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=default.target
